# kindle/kindle-dashboard.conf
# Kindle Dashboard project - reversible Upstart entry point.
description "Kindle Dashboard launcher"

start on framework_ready
stop on stopping framework

task

script
  if [ -f /mnt/us/dash-autostart.sh ]; then
    /bin/sh /mnt/us/dash-autostart.sh
  fi
end script
